# docker/crontab
# school-dashboard cron
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/local/sbin:/bin:/usr/bin:/sbin:/usr/sbin
# Set env file path so school-sync.sh finds Docker secrets
SCHOOL_DASHBOARD_ENV=/app/config/env

# Morning sync at 6:00am weekdays (data refresh only)
0 6 * * 1-5 root bash /app/sync/school-sync.sh >> /app/state/school-sync.log 2>&1
# Afternoon data refresh at 2:30pm weekdays
30 14 * * 1-5 root bash /app/sync/school-sync.sh >> /app/state/school-sync.log 2>&1

# Morning digest — 7:00am daily (7 days/week)
0 7 * * * root bash /app/sync/run-digest.sh morning
# Afternoon homework check — 3:30pm weekdays only
30 15 * * 1-5 root bash /app/sync/run-digest.sh afternoon
# Night prep — 8:30pm daily (7 days/week)
30 20 * * * root bash /app/sync/run-digest.sh night

# Weekly digests
# Friday 3pm: week in review
0 15 * * 5 root bash /app/sync/run-weekly.sh friday
# Sunday 7pm: week ahead preview
0 19 * * 0 root bash /app/sync/run-weekly.sh sunday
